# scripts/vmstation-event-wake.service
[Unit]
Description=VMStation Event-Driven Wake Monitor
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/vmstation-event-wake.sh
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
