# scripts/vmstation-event-wake.timer
[Unit]
Description=VMStation Event-Driven Wake Monitor Timer

[Timer]
OnBootSec=2min
OnUnitActiveSec=10min

[Install]
WantedBy=timers.target
